# java_console/bin/etc/systemd/system/rusefi.service
#
# sudo systemctl enable rusefi
#

[Unit]
Description=rusEFI systemd

[Service]
ExecStart=/rw/startup_root.sh

[Install]
WantedBy=multi-user.target
